# src/src_buildroot/buildroot/OVERLAY_FS/etc/udev/rules.d/85-hwclock.rules
KERNEL=="rtc0", RUN+="/sbin/hwclock --rtc=$root/$name --hctosys"
